# mysqld_exporter/contrib/mysqld_exporter.service
# -*- mode: conf -*-

[Unit]
Description=Prometheus exporter for MySQL server metrics.
Documentation=https://github.com/prometheus/mysqld_exporter
After=network.target

[Service]
EnvironmentFile=-/etc/sysconfig/mysqld_exporter
User=prometheus
ExecStart=/usr/bin/mysqld_exporter $ARGS
Restart=on-failure

[Install]
WantedBy=multi-user.target
